# image/stage-chillcheck/files/chillcheck-local-ui.service
[Unit]
Description=ChillCheck — Local Setup UI
# Start even without cloud credentials — shows connect screen
After=network.target
Wants=network.target

[Service]
Type=simple
User=chillcheck
WorkingDirectory=/opt/chillcheck
EnvironmentFile=/etc/chillcheck/.env
ExecStart=/opt/chillcheck/venv/bin/python /opt/chillcheck/local_ui/app.py
# Allow non-root user to bind to port 80
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
Restart=on-failure
RestartSec=10s
StandardOutput=append:/var/log/chillcheck/local_ui.log
StandardError=append:/var/log/chillcheck/local_ui.log

[Install]
WantedBy=multi-user.target
