# image/stage-chillcheck/files/chillcheck-firstboot.service
[Unit]
Description=ChillCheck First Boot Setup
# Run after network is up but before other ChillCheck services
After=network-online.target
Wants=network-online.target
# Only run if the flag file doesn't exist
ConditionPathExists=!/etc/chillcheck/.firstboot-complete

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/chillcheck-firstboot
StandardOutput=journal+console
StandardError=journal+console
# Give it up to 15 minutes (npm install can be slow)
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
