# image/stage-chillcheck/files/chillcheck-subscriber.service
# ============================================================
# chillcheck-subscriber.service
# ============================================================
[Unit]
Description=ChillCheck — MQTT Subscriber & Alerting
After=network.target mosquitto.service zigbee2mqtt.service
Requires=mosquitto.service

[Service]
Type=simple
User=chillcheck
WorkingDirectory=/opt/chillcheck
EnvironmentFile=/etc/chillcheck/.env
ExecStart=/opt/chillcheck/venv/bin/python /opt/chillcheck/subscriber/main.py
Restart=on-failure
RestartSec=15s
StandardOutput=append:/var/log/chillcheck/subscriber.log
StandardError=append:/var/log/chillcheck/subscriber.log

[Install]
WantedBy=multi-user.target
